# deploy/openfand.service
[Unit]
Description=OpenFAN Controller Server
Documentation=https://github.com/graelo/OpenFanController
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=openfan
Group=openfan
WorkingDirectory=/opt/openfan
ExecStart=/opt/openfan/bin/openfand --config /etc/openfan/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openfand

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/openfan /var/log/openfan /etc/openfan
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictNamespaces=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Device access for hardware communication
DeviceAllow=/dev/ttyUSB* rw
DeviceAllow=/dev/ttyACM* rw
DevicePolicy=closed

# Environment
Environment=RUST_LOG=info
Environment=OPENFAN_SERVER_CONFIG=/etc/openfan/config.toml

[Install]
WantedBy=multi-user.target
Alias=openfand.service
